# ln-faucet.service
[Unit]
Description=Bitcoin Lightning Faucet
After=network.target lnbits.service
Wants=lnbits.service

[Service]
Type=simple
User=felipe
Group=felipe
WorkingDirectory=/home/felipe/Bots/BTCfaucet
EnvironmentFile=/home/felipe/Bots/BTCfaucet/.env

# Aguarda LNbits + Phoenixd estabilizarem (ajuste conforme seu setup)
ExecStartPre=/bin/sleep 240

ExecStart=/home/felipe/Bots/BTCfaucet/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8420 \
    --workers 1 \
    --log-level info \
    --access-log

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ln-faucet

# Limites de recursos
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
